# _p2w_r619ac-128m/target-arm_cortex-a7+neon-vfpv4_musl_eabi/root.orig-ipq40xx/usr/lib/opkg/info/luci-app-adbyby-plus.list
/usr/share/adbyby/adbybyupdate.sh
/usr/lib/lua/luci/view/adbyby/refresh.htm
/usr/share/adbyby/dnsmasq.adblock
/usr/lib/lua/luci/model/cbi/adbyby/black.lua
/usr/share/adbyby/blockip.conf
/usr/share/adbyby/adblock.sh
/usr/share/adbyby/adhost.conf
/etc/init.d/adbyby
/usr/share/adbyby/adbybyfirst.sh
/usr/lib/lua/luci/model/cbi/adbyby/advanced.lua
/usr/lib/lua/luci/model/cbi/adbyby/base.lua
/usr/lib/lua/luci/model/cbi/adbyby/user.lua
/etc/config/adbyby
/usr/share/adbyby/rule-count
/usr/share/adbyby/rules.txt
/lib/upgrade/keep.d/adbyby
/usr/share/adbyby/ad-update
/usr/share/adbyby/admem.sh
/usr/lib/lua/luci/view/adbyby/adbyby_status.htm
/etc/uci-defaults/luci-adbyby
/usr/lib/lua/luci/model/cbi/adbyby/esc.lua
/usr/share/adbyby/dnsmasq.esc
/usr/share/adbyby/firewall.include
/usr/lib/lua/luci/model/cbi/adbyby/help.lua
/usr/lib/lua/luci/model/cbi/adbyby/block.lua
/usr/share/adbyby/adesc.conf
/usr/share/adbyby/adblack.conf
/usr/share/adbyby/rule-update
/usr/lib/lua/luci/controller/adbyby.lua
/usr/share/adbyby/adupdate.sh
/usr/share/adbyby/dnsmasq.ads
