# _p2w_r619ac-128m/target-arm_cortex-a7+neon-vfpv4_musl_eabi/root-ipq40xx/usr/lib/opkg/info/kmod-shortcut-fe.list
/lib/modules/5.4.158/shortcut-fe-ipv6.ko
/usr/bin/sfe_dump
/etc/modules.d/09-shortcut-fe
/lib/modules/5.4.158/shortcut-fe.ko
/etc/init.d/shortcut-fe
