# _p2w_r619ac-128m/target-arm_cortex-a7+neon-vfpv4_musl_eabi/root-ipq40xx/usr/lib/opkg/info/luci-app-turboacc.list
/etc/uci-defaults/luci-turboacc
/etc/init.d/turboacc
/usr/lib/lua/luci/model/cbi/turboacc.lua
/etc/init.d/turboacc_fss
/etc/config/turboacc
/usr/share/dnscache/dnscache-while.sh
/usr/lib/lua/luci/controller/turboacc.lua
/usr/lib/lua/luci/view/turboacc/turboacc_status.htm
