# _p2w_r619ac-128m/target-arm_cortex-a7+neon-vfpv4_musl_eabi/root-ipq40xx/usr/lib/opkg/info/luci-app-arpbind.list
/usr/lib/lua/luci/controller/arpbind.lua
/usr/lib/lua/luci/model/cbi/arpbind.lua
/etc/init.d/arpbind
/etc/uci-defaults/luci-arpbind
/usr/share/rpcd/acl.d/luci-app-arpbind.json
/etc/hotplug.d/iface/50-arpbind
