# _p2w_r619ac-128m/target-arm_cortex-a7+neon-vfpv4_musl_eabi/root-ipq40xx/usr/lib/opkg/info/luci-app-guest-wifi.list
/etc/init.d/guest-wifi
/etc/config/guest-wifi
/etc/uci-defaults/luci-guest-wifi
/usr/lib/lua/luci/controller/guest-wifi.lua
/usr/lib/lua/luci/model/cbi/guest-wifi.lua
