# _p2w_r619ac-128m/target-arm_cortex-a7+neon-vfpv4_musl_eabi/root-ipq40xx/usr/lib/opkg/info/luci-app-autoreboot.list
/usr/lib/lua/luci/controller/autoreboot.lua
/usr/lib/lua/luci/model/cbi/autoreboot.lua
/etc/init.d/autoreboot
/etc/config/autoreboot
