# _p2w_r619ac-128m/target-arm_cortex-a7+neon-vfpv4_musl_eabi/root-ipq40xx/usr/lib/opkg/info/luci-app-ipsec-vpnd.list
/etc/init.d/ipsec
/usr/lib/lua/luci/view/ipsec/ipsec_status.htm
/etc/uci-defaults/luci-ipsec
/etc/config/ipsec
/usr/lib/lua/luci/controller/ipsec-server.lua
/usr/lib/lua/luci/model/cbi/ipsec-server/ipsec-server.lua
/etc/ipsec.include
