# _p2w_r619ac-128m/target-arm_cortex-a7+neon-vfpv4_musl_eabi/root-ipq40xx/usr/lib/opkg/info/nlbwmon.list
/usr/sbin/nlbw
/etc/hotplug.d/iface/30-nlbwmon
/etc/config/nlbwmon
/usr/share/nlbwmon/protocols
/etc/init.d/nlbwmon
/usr/sbin/nlbwmon
